# udev/99-yurex.rules
ACTION=="add", BUS=="usb", SYSFS{idVendor}=="0c45", SYSFS{idProduct}=="1010", \
    RUN+="usb_yurex '%b:1.0'"
